# skills/linkedin_engagement_v1/systemd/linkedin-daily-engagement.service
[Unit]
Description=LinkedIn Daily Engagement - social agent
Wants=linkedin-daily-engagement.timer

[Service]
Type=oneshot
ExecStart=/home/freskhy/.openclaw/workspace/skills/linkedin_engagement_v1/run_engagement.sh
User=freskhy
WorkingDirectory=/home/freskhy/.openclaw/workspace-social
Environment="HOME=/home/freskhy"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/freskhy/.npm-global/bin"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
