# systemd/openclaw-daily-improvement.service
[Unit]
Description=OpenClaw Daily Improvement Task
Wants=openclaw-daily-improvement.timer

[Service]
Type=oneshot
ExecStart=/home/freskhy/projects/openclaw-workspace/scripts/daily_improvement.sh
User=freskhy
WorkingDirectory=/home/freskhy/projects/openclaw-workspace

[Install]
WantedBy=multi-user.target
