# skills/news_digest_v1/systemd/news-daily-digest.timer
[Unit]
Description=Daily News Digest Timer (8 AM)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
